# q3/q3.cron
* * * * * /home/capitan/work/questions-vf/q3/ps2.sh
